# dozzle-llm-watch.service
[Unit]
Description=Dozzle LLM Log Analyzer
After=network.target

[Service]
Type=simple
User=theanswer
WorkingDirectory=/opt/dozzle-llm-watch

ExecStart=/opt/dozzle-llm-watch/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8088 --timeout-keep-alive 65

Restart=always
RestartSec=5

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
